# service/transmission/run
#!/bin/sh

exec /sbin/setuser abc transmission-daemon -f -g /config  --no-incomplete-dir --no-auth -w /download --no-portmap --dht --global-seedratio 1.0 -ep
