# service/transmission/run
#!/bin/sh

exec /sbin/setuser abc transmission-daemon -g /config  --no-incomplete-dir --no-auth -w /download